# workstation_setups/OL9_ZFS_Server/ipmi-fan-override.service
[Unit]
Description=Supermicro IPMI Fan Override
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ipmi-fan-override.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
